# puppet/modules/riemann/templates/dashboard/upstart.conf
#!upstart

description "Riemann Web GUI"
start on runlevel [23]
stop on shutdown
chdir /srv/riemann-dash/
<% bin_path = (`echo "$(gem env)"|grep DIREC|grep EXEC| awk '{print $4}'`).strip %>
exec su -c "<%= bin_path %>/riemann-dash >> /srv/riemann-dash/riemann-dash.log 2>&1" riemann

respawn
