# puppet/modules/riemann/templates/server/upstart.conf
#!upstart

description "Riemann"
start on runlevel [23]
stop on shutdown

exec sudo -u riemann /opt/riemann/bin/riemann /etc/riemann/riemann.config

respawn
